# altrpnzssr.conf
[program:altrpnzssr]
directory=/Users/ilasabusaev/ols-docker-env/sites/altrp.nz
command=node server-dist/index.js
numprocs=1
autostart=true
autorestart=true
user=ilasabusaev